# deploy/usage-analytics/cybercode-usage.service
[Unit]
Description=CyberCode anonymous usage statistics
After=network.target

[Service]
Type=simple
User=cybercode-usage
Group=cybercode-usage
WorkingDirectory=/opt/cybercode-usage
EnvironmentFile=/etc/cybercode-usage.env
Environment=PYTHONPATH=/opt/cybercode-usage/vendor
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 /opt/cybercode-usage/server.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cybercode-usage
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
